# bridge/com.james.buddykinit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.james.buddykinit</string>

  <!-- Renew the FILE-based Kerberos cache every 12 h.
       The wrapper script tries `kinit -R` first (cheap, password-free)
       and falls back to a full kinit using BUDDY_KRB_USER / BUDDY_KRB_PASS
       sourced from ~/.zshrc (chmod 600).  Neither file in the repo
       contains the secret — only the user's local shell rc does. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/bytedance/buddy/bridge/krb-renew.sh</string>
  </array>

  <key>RunAtLoad</key><true/>
  <key>StartInterval</key><integer>43200</integer>   <!-- 12 hours -->
  <key>StandardOutPath</key><string>/tmp/buddykinit.out</string>
  <key>StandardErrorPath</key><string>/tmp/buddykinit.err</string>
</dict>
</plist>
